# deploy/supervisord.conf
; ============================================================================
;  Independent Commerce — supervisord configuration
;  Runs the FastAPI backend and nginx side by side inside a single container.
;  Coolify sees only port 80 (nginx).
; ============================================================================

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; --- FastAPI backend --------------------------------------------------------
; Single worker: seeding runs once (no race on MongoDB unique indexes) and
; memory footprint stays low on small Coolify VMs. Nginx handles concurrency
; via keepalive to this worker.
;
; The command is wrapped in `sh -c 'echo ...; exec uvicorn ...'` so that
; a "BOOT:" line ALWAYS appears in the stdout log the instant supervisor
; spawns the process — even if uvicorn itself hangs during module import
; (e.g. a slow MongoDB SRV DNS lookup). This makes debugging silent
; startup hangs on Coolify possible via `supervisorctl tail backend`.
;
; `redirect_stderr=true` merges uvicorn's stderr (which is where it emits
; its own INFO/ACCESS logs by default) into the stdout channel so a
; single `supervisorctl tail backend` shows everything.
[program:backend]
command=sh -c 'echo "BOOT: launching uvicorn (server:app on 127.0.0.1:8001)"; exec python -u -m uvicorn server:app --host 127.0.0.1 --port 8001 --workers 1 --proxy-headers --forwarded-allow-ips=* --timeout-keep-alive 30 --log-level info'
directory=/app/backend
autostart=true
autorestart=true
startsecs=3
stopsignal=INT
stopwaitsecs=10
priority=10
redirect_stderr=true
stdout_logfile=/var/log/supervisor/backend.stdout.log
stderr_logfile=/var/log/supervisor/backend.stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
environment=PYTHONUNBUFFERED="1",PYTHONFAULTHANDLER="1"

; --- nginx reverse proxy + static frontend ---------------------------------
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startsecs=2
stopsignal=QUIT
stopwaitsecs=10
priority=20
stdout_logfile=/var/log/supervisor/nginx.stdout.log
stderr_logfile=/var/log/supervisor/nginx.stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
