# deploy/supervisord.conf
; ============================================================================
;  Independent Commerce — supervisord configuration
;  Runs the FastAPI backend and nginx side by side inside a single container.
;  Coolify sees only port 80 (nginx).
; ============================================================================

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; --- FastAPI backend --------------------------------------------------------
[program:backend]
command=/usr/local/bin/uvicorn server:app --host 127.0.0.1 --port 8001 --workers 2 --proxy-headers --forwarded-allow-ips=*
directory=/app/backend
autostart=true
autorestart=true
startsecs=3
stopsignal=INT
stopwaitsecs=10
priority=10
stdout_logfile=/var/log/supervisor/backend.stdout.log
stderr_logfile=/var/log/supervisor/backend.stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
environment=PYTHONUNBUFFERED=1

; --- nginx reverse proxy + static frontend ---------------------------------
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startsecs=2
stopsignal=QUIT
stopwaitsecs=10
priority=20
stdout_logfile=/var/log/supervisor/nginx.stdout.log
stderr_logfile=/var/log/supervisor/nginx.stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
